# deploy/pep-oracle-api.service
[Unit]
Description=PEP Oracle API server
After=network.target

[Service]
Type=simple
User=pep-oracle
WorkingDirectory=/opt/pep-oracle/app
ExecStart=/opt/pep-oracle/app/.venv/bin/pep-oracle-server
EnvironmentFile=/opt/pep-oracle/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
