# deploy/pep-oracle-ingest.service
[Unit]
Description=PEP Oracle episode ingestion
After=network.target
# Run the off-site backup only after a successful ingest (the corpus only changes here).
OnSuccess=pep-oracle-backup.service

[Service]
Type=oneshot
User=gregg
WorkingDirectory=/opt/pep-oracle/app
ExecStart=/opt/pep-oracle/app/.venv/bin/pep-oracle ingest --diarize --new-only
EnvironmentFile=/opt/pep-oracle/.env
